# server/supervisor_staging.conf
[program:brochachos_server]
directory=/home/hta/Programs/ultimateContacts/server/
command=sails lift --port=1333
user=hta
autorestart=true
autostart=true
