# contact-book/supervisor_staging.conf
[program:brochachos]
directory=/home/hta/Programs/ultimateContacts/contact-book/
command=ember serve --environment=staging
user=hta
autorestart=true
autostart=true
